# scripts/llmready-backend.service
[Unit]
Description=LLMReady FastAPI Backend
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=llmready
WorkingDirectory=/opt/llmready/backend
Environment="PATH=/opt/llmready/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/llmready/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target