# scripts/llmready-celery-worker.service
[Unit]
Description=LLMReady Celery Worker
After=network.target llmready-backend.service
Requires=llmready-backend.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/llmready/backend
Environment="PATH=/opt/llmready/venv/bin"

# Load environment variables from backend .env
EnvironmentFile=/opt/llmready/backend/.env

# Start Celery worker
ExecStart=/opt/llmready/venv/bin/celery -A app.core.celery_app worker \
    --loglevel=info \
    --concurrency=2 \
    --max-tasks-per-child=1000

# Restart on failure
Restart=always
RestartSec=10s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=celery-worker

[Install]
WantedBy=multi-user.target